# var/lib/dpkg/info/postgresql-common.conffiles
/etc/init.d/postgresql
/etc/logrotate.d/postgresql-common
